# scripts/vm-bootstrap/linux/systemd/gh-runner-watchdog.service
[Unit]
Description=GitHub Actions runner idle/lifetime watchdog
ConditionPathExists=/opt/gh-runner-bootstrap/watchdog.sh

[Service]
Type=oneshot
User=root
ExecStart=/opt/gh-runner-bootstrap/watchdog.sh
